# firmware/nonos_server/shairplay.service
[Unit]
Description=Shairplay
After=network.target
Wants=network.target

[Service]
Type=simple
User=atopile
WorkingDirectory=/home/atopile/shairplay
ExecStart=/usr/local/bin/shairplay
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
